# supervisord.conf
[supervisord]
nodaemon=true

[program:kestrel]
directory=/src
command=dnx . kestrel
